# deploy/systemd/strele-grid-cache.timer
[Unit]
Description=Run Strelko 1km grid cache every 5 minutes
# Zahteva delujoč Docker kontejner strele2-obcine-public-1 (service zaganja job prek docker exec).

[Timer]
# Vsakih 5 min — usklajeno z ingestom (~5 min) in on-demand API osvežitvijo (~2 min).
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target

